# etc/catkin/profile.d/15.rosbash.bash
# generated from rosbash/env-hooks/15.rosbash.bash.em

if [ -z "$CATKIN_ENV_HOOK_WORKSPACE" ]; then
  CATKIN_ENV_HOOK_WORKSPACE="/home/scmbuild/workspaces_cluster/baidu.car-os.ros/baidu/car-os/install_isolated"
fi
. "$CATKIN_ENV_HOOK_WORKSPACE/share/rosbash/rosbash"
